# roles/security_veracrypt/templates/vc.service.j2
[Install]
WantedBy=default.target

[Unit]
Description=VeraCrypt (%i) - mount user container

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/vc open {{ item.path }} {{ item.mount }} vc_{{ item.name }}
ExecStop=/usr/local/bin/vc close {{ item.mount }} vc_{{ item.name }}

[Install]
WantedBy=default.target
